# packages/Quackoubert/quackoubert-0.0.0.tar.gz/quackoubert-0.0.0/crontab
@reboot (cd /home/marceau/scripts; bash /home/marceau/scripts/start_europarser.sh) && echo $(date -u) "Europarser reboot start sucessful" >> /home/marceau/scripts/log_whole_start || echo $(date -u) "Europarser reboot start failed" >> /home/marceau/scripts/log_whole_start
@reboot (cd /home/marceau/scripts; bash /home/marceau/scripts/start_TXM4TXM.sh) && echo $(date -u) "TXM4TXM reboot start sucessful" >> /home/marceau/scripts/log_whole_start || echo $(date -u) "TXM4TXM reboot start failed" >> /home/marceau/scripts/log_whole_start
@reboot (cd /home/marceau/scripts; bash /home/marceau/scripts/start_EuroToDB.sh) && echo $(date -u) "EuroToDB reboot start sucessful" >> /home/marceau/scripts/log_whole_start || echo $(date -u) "EuroToDB reboot start failed" >> /home/marceau/scripts/log_whole_start
