# packages/inmanta-module-mysql/inmanta_module_mysql-0.6.15-py3-none-any.whl/inmanta_plugins/mysql/files/mysql_create_db.service
[Unit]
Description=MySQL database create script
After=mysqld.target

[Service]
Type=simple
User=mysql
Group=mysql

ExecStart=/usr/bin/mysql_create_db

[Install]
WantedBy=multi-user.target

